# scripts/nixos-install.service
[Unit]
Description=Automated NixOS install (ZFS/Disko)
Wants=basic.target
After=basic.target
ConditionPathExists=/root/installer/install.sh

[Service]
Type=oneshot
Environment=NIXOS_SNAPSHOT_GATE=0
Environment=NIXOS_SNAPSHOT_WAIT_SEC=300
ExecStart=/root/installer/install.sh
StandardOutput=journal+console
StandardError=journal+console
TimeoutStartSec=0
Restart=on-failure
RestartSec=15

[Install]
WantedBy=multi-user.target
WantedBy=default.target
